# scripts/quant-lab.service
[Unit]
Description=Quant Lab Runtime API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__PROJECT_DIR__
ExecStart=__PROJECT_DIR__/.venv/bin/python -m quant_lab service-api --config config/settings.yaml
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
